# roles/meraki_exporter/templates/meraki_exporter.service.j2
# =============================================================================
# meraki_exporter.service — systemd unit for Meraki Prometheus exporter
# =============================================================================
# Generated by Ansible — do not edit manually.
# Source: roles/meraki_exporter/templates/meraki_exporter.service.j2
# =============================================================================

[Unit]
Description=Meraki Prometheus Exporter
Documentation=https://github.com/ptchuitio84/homelabinfra-iac-ansible
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ meraki_exporter_user }}
Group={{ meraki_exporter_group }}
ExecStart={{ meraki_exporter_venv }}/bin/python3 {{ meraki_exporter_script }}
Restart=on-failure
RestartSec=15s
StandardOutput=append:{{ meraki_exporter_log_dir }}/meraki_exporter.log
StandardError=append:{{ meraki_exporter_log_dir }}/meraki_exporter.log

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ meraki_exporter_log_dir }}

[Install]
WantedBy=multi-user.target
