# roles/unpoller/templates/unpoller.service.j2
# =============================================================================
# unpoller.service — systemd unit for UnPoller
# =============================================================================
# Generated by Ansible — do not edit manually.
# Source: roles/unpoller/templates/unpoller.service.j2
# =============================================================================

[Unit]
Description=UnPoller — UniFi Prometheus Exporter
Documentation=https://unpoller.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ unpoller_user }}
Group={{ unpoller_group }}
ExecStart={{ unpoller_install_dir }}/unpoller --config {{ unpoller_config_file }}
Restart=on-failure
RestartSec=15s
StandardOutput=append:{{ unpoller_log_dir }}/unpoller.log
StandardError=append:{{ unpoller_log_dir }}/unpoller.log

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths={{ unpoller_log_dir }}

[Install]
WantedBy=multi-user.target
