# roles/loki/templates/loki.service.j2
# =============================================================================
# /etc/systemd/system/loki.service — managed by Ansible (roles/loki)
# =============================================================================

[Unit]
Description=Grafana Loki
After=network.target

[Service]
User={{ loki_user }}
Group={{ loki_group }}
Type=simple
ExecStart={{ loki_install_dir }}/loki -config.file={{ loki_config_dir }}/loki.yml
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
