# roles/prometheus/templates/prometheus.service.j2
# =============================================================================
# prometheus.service — Prometheus systemd Unit File
# =============================================================================
# Generated by Ansible — do not edit manually.
# Source: roles/prometheus/templates/prometheus.service.j2
# =============================================================================

[Unit]
Description=Prometheus Monitoring System
Documentation=https://prometheus.io/docs/introduction/overview/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ prometheus_user }}
Group={{ prometheus_group }}
ExecStart={{ prometheus_install_dir }}/prometheus \
    --config.file={{ prometheus_config_dir }}/prometheus.yml \
    --storage.tsdb.path={{ prometheus_data_dir }} \
    --storage.tsdb.retention.time={{ prometheus_retention }} \
    --web.listen-address=0.0.0.0:{{ prometheus_port }} \
    --web.enable-lifecycle

# Restart automatically if Prometheus crashes.
Restart=on-failure
RestartSec=5s

# Security hardening — Prometheus does not need root.
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
